# palsrv.service
[Unit]
Description=Palworld Dedicated Server

[Service]
ExecStart=/srv/palworld/run.sh server
Restart=on-failure

[Install]
WantedBy=default.target
